# package/alpine/vane.initd
#!/sbin/openrc-run

name="vane"
description="Vane - Flow-based reverse proxy"

command="/usr/bin/vane"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
directory="/etc/vane"

depend() {
    need net
    after firewall
}

start_pre() {
    if [ -f "/etc/vane/.env" ]; then
        while IFS= read -r line; do
            case "$line" in
                \#*|"") continue ;;
                *=*) export "$line" ;;
            esac
        done < "/etc/vane/.env"
    fi
}
