# package/alpine/vane.post-install
#!/bin/sh

echo "Vane has been installed."
echo "Configuration directory: /etc/vane"
echo "To start: rc-service vane start"
echo "To enable on boot: rc-update add vane"
